# output/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:redis]
command=redis-server --appendonly yes --maxmemory 128mb --maxmemory-policy allkeys-lru
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-err.log

[program:agent-po]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/agents/po
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/po.log
stderr_logfile=/var/log/supervisor/po-err.log

[program:agent-pm]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8001
directory=/app/agents/pm
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/pm.log
stderr_logfile=/var/log/supervisor/pm-err.log

[program:agent-ba]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8002
directory=/app/agents/ba
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/ba.log
stderr_logfile=/var/log/supervisor/ba-err.log

[program:agent-sa]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8003
directory=/app/agents/solution_architect
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/sa.log
stderr_logfile=/var/log/supervisor/sa-err.log

[program:agent-frontend-dev]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8004
directory=/app/agents/frontend_dev
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/frontend-dev.log
stderr_logfile=/var/log/supervisor/frontend-dev-err.log

[program:agent-backend-dev]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8005
directory=/app/agents/backend_dev
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/backend-dev.log
stderr_logfile=/var/log/supervisor/backend-dev-err.log

[program:agent-uiux-reviewer]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8006
directory=/app/agents/uiux_reviewer
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/uiux-reviewer.log
stderr_logfile=/var/log/supervisor/uiux-reviewer-err.log

[program:agent-security-reviewer]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8007
directory=/app/agents/security_reviewer
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/security-reviewer.log
stderr_logfile=/var/log/supervisor/security-reviewer-err.log

[program:agent-qa]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8008
directory=/app/agents/qa
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/qa.log
stderr_logfile=/var/log/supervisor/qa-err.log

[program:agent-devops]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8009
directory=/app/agents/devops
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/devops.log
stderr_logfile=/var/log/supervisor/devops-err.log

[program:agent-techlead]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8010
directory=/app/agents/techlead
autostart=true
autorestart=true
environment=PYTHONPATH="/app"
stdout_logfile=/var/log/supervisor/techlead.log
stderr_logfile=/var/log/supervisor/techlead-err.log

[program:web]
command=node server.js
directory=/app/web/standalone
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="7860",HOSTNAME="0.0.0.0",NEXT_PUBLIC_API_URL="http://localhost:7860/api",PO_URL="http://localhost:8000",PM_URL="http://localhost:8001",BA_URL="http://localhost:8002",SA_URL="http://localhost:8003",FRONTEND_DEV_URL="http://localhost:8004",BACKEND_DEV_URL="http://localhost:8005",UIUX_REVIEWER_URL="http://localhost:8006",SECURITY_REVIEWER_URL="http://localhost:8007",QA_URL="http://localhost:8008",DEVOPS_URL="http://localhost:8009",TECHLEAD_URL="http://localhost:8010",REDIS_URL="redis://localhost:6379/0"
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web-err.log
